# custom_pigen/stage2/04-install-requirements/files/99-wlan0.rules
SUBSYSTEM=="net", KERNEL=="wlan0", RUN+="/bin/systemctl restart ap0-setup.service"
