# custom_pigen/stage2/04-install-requirements/files/ap0-setup.service
[Unit]
Description=Create ap0 interface
After=sys-subsystem-net-devices-wlan0.device

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/sbin/iw phy phy0 interface add ap0 type __ap
ExecStartPost=/usr/bin/ip link set dev ap0 address 00:11:22:33:44:55
ExecStartPost=/sbin/ifconfig ap0 up
ExecStop=/sbin/iw dev ap0 del

[Install]
WantedBy=multi-user.target
