# custom_pigen/stage2/04-install-requirements/files/simpleaq.service
[Unit]
Description=SimpleAQ Recorder

[Install]
WantedBy=default.target

[Service]
Restart=on-failure
ExecStart=python /simpleaq/simpleaq.py --env /etc/environment

